# backup/root.crontab
*/5 * * * * /usr/bin/flock -xn /root/.backup.lock -c /root/backup.sh >> /var/log/backup.log 2>&1
